# daemon/appimage-manager.service.in
[Unit]
Description=AppImage Manager daemon
After=graphical-session.target

[Service]
Type=simple
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/appimage-manager-daemon
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
